# files/openvpn-rc.d
#!/bin/sh

daemon="/usr/local/sbin/openvpn"

. /etc/rc.d/rc.subr

rc_cmd $1
